# rpi-gpiod.service
[Unit]
Description="Raspberry Pi GPIO daemon"

[Service]
ExecStart=/usr/bin/rpi-gpiod

[Install]
WantedBy=multi-user.target
